# deploy/systemd/robot-hand.service
[Unit]
# =============================================================================
# BrazoRoboticoV2 — Flask + WebSocket
#
# Servidor principal del brazo robótico. Sirve la API REST, WebSocket para
# tiempo real, y el frontend estático (React build).
#
# Depende de socat-mega.service para el bridge con el Mega por USB.
# =============================================================================
Description=BrazoRoboticoV2 — Flask + WebSocket
After=network-online.target socat-mega.service
Wants=socat-mega.service

[Service]
Type=simple
User=arduino
WorkingDirectory=/home/arduino/brazorobotico/backend
ExecStart=/usr/bin/python3 server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
