# deploy/systemd/socat-mega.service
[Unit]
# =============================================================================
# SOCAT bridge: TCP:7500 ↔ Mega 2560 por USB OTG
#
# El Mega 2560 se conecta al UNO Q mediante USB OTG y aparece como /dev/ttyACM0.
# SOCAT crea un bridge TCP para que Flask pueda hablar con el Mega
# a través del puerto serie virtual.
#
# Flujo: Flask → TCP:7500 → SOCAT → /dev/ttyACM0 → Mega 2560 → PWM → Servos
# =============================================================================
Description=SOCAT bridge: TCP:7500 ↔ Mega 2560 (USB OTG)
After=multi-user.target
Requires=dev-ttyACM0.device

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/socat file:/dev/ttyACM0,raw,echo=0,b115200 tcp-listen:7500,reuseaddr,fork
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
